# scripts/spec-upstart.conf
#!upstart
#
description "spec - events manager built with node.js"
author      "Robert Christensen <rchristensen@wesleyan.edu>"
 
start on startup
stop on shutdown
 
# We do this so that upstart will report the PID of the node.js process, rather than forever's
expect fork
 
# Environmental variables for node.js/pulleffect
env NODE_BIN_DIR="/usr/bin"
env NODE_PATH="/usr/lib/node_modules"
env APPLICATION_DIRECTORY="/var/www/spec"
env APPLICATION_START="spec-server.js"
env LOG="/var/log/spec.log"
 
script
    # Add the node executables to the path
    PATH=$NODE_BIN_DIR:$PATH
    # The minUptime and spinSleepTime settings stop Forever from thrashing if
    # the application fails immediately on launch. This is generally necessary to
    # avoid loading development servers to the point of failure every time 
    # someone makes an error in application initialization code, or bringing down
    # production servers the same way if a database or other critical service
    # suddenly becomes inaccessible.
    exec forever --sourceDir $APPLICATION_DIRECTORY -a -l $LOG \
         --minUptime 5000 --spinSleepTime 2000 start $APPLICATION_START
end script
 
pre-stop script
    # Add the node executables to the path.
    PATH=$NODE_BIN_DIR:$PATH
    # Here we're using the pre-stop script to stop the Node.js application
    # process so that Forever is given a chance to do its thing and tidy up
    # its data. Note that doing it this way means that each application that
    # runs under Forever must have a different start file name, regardless of
    # which directory it is in.
    exec forever stop $APPLICATION_START >> $LOG
end script